# scripts/com.hkjc.dailypipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hkjc.dailypipeline</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/fatlung/.openclaw/workspace-main/hkjc_project/daily_pipeline.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/fatlung/.openclaw/workspace-main/hkjc_project</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/fatlung/.openclaw/workspace-main/hkjc_project</string>
        <key>MONGODB_DB_NAME</key>
        <string>hkjc_racing_dev</string>
    </dict>
    
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/fatlung/.openclaw/workspace-main/hkjc_project/logs/pipeline/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/fatlung/.openclaw/workspace-main/hkjc_project/logs/pipeline/stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
